# docker-gen/dist/docker-gen.service
[Unit]
Description=docker-gen

[Service]
PIDFile=/var/run/docker-gen
User=root
Group=root
WorkingDirectory=/opt/docker-gen
ExectStart=/opt/docker-gen

[Install]
WantedBy=multi-user.target
