# jupyter/etc/profile.d/basher.sh
export BASHER_ROOT=$HOME/.config/basher
